# roles/access-point/files/99-wlan1.rules
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlan0", RUN+="/usr/sbin/iw dev wlan0 interface add wlan1 type __ap"